# service/php/run
#!/bin/sh -e
exec 2>&1
exec php-fpm81 -F