# service/nginx/run
#!/bin/sh -e
exec 2>&1
exec env - PATH=$PATH nginx -g 'daemon off;'
exec php-fpm81 -F